# examples/network-manager.service
[Unit]
Description=Network Manager - Automatische WiFi/LAN Verwaltung
After=NetworkManager.service
Wants=NetworkManager.service

[Service]
Type=simple
ExecStart=/usr/local/bin/network-manager --daemon
Restart=always
RestartSec=5
User=%i

[Install]
WantedBy=multi-user.target
